# restart.if.no.net.timer
[Unit]
Description=Restart the router if network is not available

[Timer]
OnCalendar=*:0/29
Persistent=false

[Install]
WantedBy=timers.target
