# restart.if.no.net.service
[Unit]
Description=Restart the router if network is not available

[Service]
ExecStart="/usr/local/bin/restart.if.no.net.sh"
